# stage3/searchwing-antennabox/rootfs/etc/rc.local
#!/bin/sh -e

# unblock wlan on every boot just in case
rfkill unblock wlan
